# systemctl/et.service
[Unit]
Description=Eternal Terminal
After=syslog.target network.target

[Service]
Type=simple
Restart=on-failure
ExecStart=/usr/bin/etserver --cfgfile=/etc/et.cfg --logtostdout

[Install]
WantedBy=multi-user.target
